# supervisord.conf
[supervisord]
nodaemon=true

[program:aerospike-amc]
command=/bin/sh -c "/opt/amc/bin/start_amc_cron.sh"
